# node/etc/systemd/system/backup.service.d/00-triggers.conf
[Unit]
Before=shutdown.target
[Install]
WantedBy=shutdown.target
