# node/etc/systemd/system/backup.service.d/00-triggers.conf
mkdir -vp /etc/systemd/system/backup.service.d
cat << 'EOF' > /etc/systemd/system/backup.service.d/00-triggers.conf
[Unit]
Before=shutdown.target
[Install]
WantedBy=shutdown.target
EOF