# node/usr/local/lib/systemd/system/backup.timer
[Unit]
Description=Timer for backup service

[Timer]
Unit=backup.service
OnCalendar=00:00
AccuracySec=15s

[Install]
WantedBy=timers.target
