# node/usr/local/lib/systemd/system/backup.timer
[Unit]
Description=Timer for backup service

[Timer]
Unit=backup.service
WakeSystem=on

[Install]
WantedBy=timers.target